# debian/postrm
#!/bin/bash
set -e

update-rc.d kiwix remove